# tests/modules/services/taskwarrior-sync/custom-frequency-expected.timer
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=hourly
Unit=taskwarrior-sync.service

[Unit]
Description=Taskwarrior periodic sync
