# tests/modules/services/taskwarrior-sync/custom-package-expected.service
[Service]
CPUSchedulingPolicy=idle
ExecStart=@taskwarrior@/bin/task synchronize
IOSchedulingClass=idle

[Unit]
Description=Taskwarrior sync
