# tests/modules/services/tahoe-lafs/basic-service-expected.service
[Unit]
Description=Tahoe-LAFS

[Service]
ExecStart=@tahoelafs@/bin/tahoe run -C %h/.tahoe