# tests/modules/services/lorri/basic-socket-expected.socket
[Install]
WantedBy=sockets.target

[Socket]
ListenStream=%t/lorri/daemon.socket
RuntimeDirectory=lorri

[Unit]
Description=Socket for lorri build daemon
