# tests/modules/services/systembus-notify/basic-service-expected.service
[Unit]
Description=systembus-notify daemon

[Service]
ExecStart=@systembus-notify@/bin/systembus-notify

[Install]
WantedBy=graphical-session.target