# tests/modules/services/mbsync/basic-service-expected.service
[Service]
ExecStart=@isync@/bin/mbsync --all --verbose
Type=oneshot

[Unit]
Description=mbsync mailbox synchronization