# tests/modules/services/mbsync/basic-service-expected.timer
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=*:0/5
Unit=mbsync.service

[Unit]
Description=mbsync mailbox synchronization