# tests/modules/services/mbsync/custom-frequency-expected.timer
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=*:0/10
Unit=mbsync.service

[Unit]
Description=mbsync mailbox synchronization