# tests/modules/services/mbsync/custom-frequency-expected.service
[Service]
ExecStart=@isync@/bin/mbsync --all
Type=oneshot

[Unit]
Description=mbsync mailbox synchronization