# tests/modules/services/ludusavi/basic-service-expected.service
[Service]
ExecStart=@ludusavi@/bin/ludusavi backup --force
Type=oneshot

[Unit]
Description=Run a game save backup with Ludusavi
