# tests/modules/services/ludusavi/basic-service-expected.timer
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=daily

[Unit]
Description=Run a game save backup with Ludusavi
