# tests/modules/services/vdirsyncer/custom-config-expected.timer
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=daily
Unit=vdirsyncer.service

[Unit]
Description=vdirsyncer calendar&contacts synchronization
