# tests/modules/services/grobi/basic-configuration-expected.service
[Install]
WantedBy=graphical-session.target

[Service]
Environment=PATH=/nix/store/00000000000000000000000000000000-xrandr/bin:@bash-interactive@/bin
ExecStart=@grobi@/bin/grobi watch -v
Restart=always
RestartSec=2s
Type=simple

[Unit]
After=graphical-session.target
Description=grobi display auto config daemon
PartOf=graphical-session.target
