# tests/modules/services/random-background/random-background-expected.timer
[Install]
WantedBy=timers.target

[Timer]
OnUnitActiveSec=30m

[Unit]
Description=Set random desktop background using feh
